# cron.d/cgiproxy
13 2 * * * root /app/cgiproxy/bin/nph-proxy.cgi purge-db
